# templates/etc/cron.d/xtrabackup.j2
# {{ ansible_managed }}

# daily full backup, purging old >= 1 day
5 1 * * * root xtrabackup_helper -r 1 >/dev/null

# incremental backup
23 2-23 * * * root xtrabackup_helper -i >/dev/null
